# agent.service
[Unit]
Description=RPI Personal AI Agent
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=tomas
WorkingDirectory=/home/tomas/agent
EnvironmentFile=/home/tomas/agent/.env
ExecStart=/home/tomas/agent/.venv/bin/python main.py telegram
Restart=on-failure
RestartSec=10
StandardOutput=append:/home/tomas/agent/agent.log
StandardError=append:/home/tomas/agent/agent.log

[Install]
WantedBy=multi-user.target

# Instalace:
#   sudo cp /home/tomas/agent/agent.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable agent
#   sudo systemctl start agent
#
# Správa:
#   sudo systemctl status agent
#   sudo systemctl restart agent
#   sudo systemctl stop agent
#   journalctl -u agent -f
